# crontab
# Daily chart sync at 06:00 in the container timezone (Europe/Amsterdam in docker-compose).
0 6 * * * cd /app && node src/sync.js
